# root/etc/services.d/xupnpd/run
#!/usr/bin/with-contenv bash

exec \
s6-setuidgid abc /etc/xupnpd/xupnpd
